# profile.d/30-unset-gtk-theme.sh
if [ x"$UNSET_GTK_THEME" != x"" ]; then
  unset GTK_THEME
  unset UNSET_GTK_THEME
fi
